# backend/docker/dev/s6-rc.d/authbackend/run
#!/command/execlineb -P
gow run -race ./cmd/authbackend --dev --log debug
