# backend/docker/dev/s6-rc.d/api/run
#!/command/execlineb -P
gow run -race ./cmd/api --dev --log debug
