# backend/docker/dev/s6-rc.d/events/run
#!/command/execlineb -P
gow run -race ./cmd/events --dev --log debug
