# backend/docker/dev/s6-rc.d/worker/run
#!/command/execlineb -P
gow run -race ./cmd/worker --dev --log debug
